# infra/systemd/daymind-api.service
[Unit]
Description=DayMind FastAPI Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=daymind
Group=daymind
EnvironmentFile=/etc/daymind/daymind.env
WorkingDirectory=/opt/daymind
ExecStart=/opt/daymind/.venv/bin/uvicorn src.api.app:app --host 0.0.0.0 --port 8000 --workers 2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
