# infra/systemd/daymind-fava.service
[Unit]
Description=DayMind Fava Ledger UI
After=network.target

[Service]
User=daymind
Group=daymind
WorkingDirectory=/opt/daymind
EnvironmentFile=/etc/default/daymind
Environment=PYTHONPATH=/opt/daymind
ExecStart=/opt/daymind/venv/bin/fava --host 127.0.0.1 --port ${FAVA_PORT} /opt/daymind/runtime/ledger.beancount
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
