# infra/systemd/daymind-fava.service
[Unit]
Description=DayMind Fava Ledger UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/daymind
EnvironmentFile=/etc/default/daymind
Environment=PYTHONPATH=/opt/daymind
ExecStartPre=/usr/bin/test -x /opt/daymind/venv/bin/fava
ExecStartPre=/usr/bin/test -d /opt/daymind
User=daymind
Group=daymind
ExecStart=/opt/daymind/scripts/start_fava.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
