# infra/systemd/daymind-fava.service
[Unit]
Description=DayMind Fava Ledger UI
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/daymind
User=daymind
Group=daymind
WorkingDirectory=/opt/daymind
ExecStart=/opt/daymind/scripts/start_fava.sh
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
